# ansible/systemd-docker/templates/ansible-pull.service.j2
[Unit]
Description=Ansible Pull Service

[Service]
Type=oneshot

ExecStart=/usr/bin/ansible-pull \
    -U {{ url }} \
    --full \
    {{ playbook }}

PrivateTmp=true
