# ansible/systemd-docker/templates/ansible-pull.timer.j2
[Unit]
Description=Ansible Pull Timer

[Timer]
OnCalendar=*:0/{{ minutes }}
Persistent=true

[Install]
WantedBy=timers.target
